# net/fika-opennds/files/buildroot/etc/hotplug.d/iface/96-fika-opennds
#!/bin/sh

. /lib/functions.sh

/etc/init.d/fika-opennds enabled && [ "$ACTION" = "ifup" ] && {
    if [ ${INTERFACE} = 'guest' ]; then
        /etc/init.d/fika-opennds restart
    fi
}
